# joystick.service
[Unit]
Description=CyberDog Joystick

[Service]
User=root
WorkingDirectory=/home/mi/Desktop/cyberdog-joystick
ExecStart=/usr/bin/python3 /home/mi/Desktop/cyberdog-joystick/joystick.py
Restart=always

[Install]
WantedBy=multi-user.target